# archiso/airootfs/etc/skel/.config/systemd/user/plasma-login-kwin_wayland.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kwin_wayland --no-lockscreen --no-global-shortcuts --no-kactivities --locale1
